# supervisord.conf
[supervisord]

nodaemon=true

[program:mongodb]

command=/start.sh mongod --dbpath /data/db
